# systemd/eoip-helper.service
[Unit]
Description=EoIP-rs Privileged Helper
Documentation=https://github.com/manawenuz/EoIP-rs
Before=eoip-rs.service

[Service]
Type=simple
ExecStart=/usr/local/bin/eoip-helper --mode persist --socket-path /run/eoip-rs/helper.sock
ExecStartPre=/bin/mkdir -p /run/eoip-rs
ExecStartPre=/bin/rm -f /run/eoip-rs/helper.sock
RuntimeDirectory=eoip-rs
RuntimeDirectoryMode=0755
Restart=on-failure
RestartSec=3

# Needs root for TAP creation and raw sockets
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
